# restart_whatsapp_services.timer
[Unit]
Description=Restart WhatsApp services every 2 hours

[Timer]
OnCalendar=*-*-* */2:00:00
Persistent=true

[Install]
WantedBy=timers.target
